# gentoo_remapper.dactyl_r
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="CyLei Dactyl_74_R remapper"
command=/usr/bin/gentoo_remapper.py
command_args="'CyLei Dactyl_74_R'"
command_background="true"
pidfile="/run/gentoo_remapper.dactyl_r.pid"
